# supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
user=root

[program:bigphish]
command=/app/venv/bin/python3 /app/bigphish.py
directory=/app
user=bigphish
autostart=true
autorestart=true
startsecs=10
startretries=3
redirect_stderr=true
stdout_logfile=/app/logs/bigphish.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/app/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:dashboard]
command=/app/venv/bin/python3 /app/metrics_dashboard.py
directory=/app
user=bigphish
autostart=true
autorestart=true
startsecs=5
startretries=3
redirect_stderr=true
stdout_logfile=/app/logs/dashboard.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:prometheus]
command=/usr/bin/prometheus --config.file=/app/prometheus.yml --storage.tsdb.path=/prometheus
directory=/app
user=bigphish
autostart=true
autorestart=true
startsecs=10
redirect_stderr=true
stdout_logfile=/app/logs/prometheus.log

[program:grafana]
command=/usr/sbin/grafana-server --homepath=/usr/share/grafana
directory=/app
user=bigphish
autostart=true
autorestart=true
startsecs=10
redirect_stderr=true
stdout_logfile=/app/logs/grafana.log

[group:bigphish]
programs=bigphish,dashboard,prometheus,grafana
priority=999